# packer/files/itamae/roles/ci/files/etc/systemd/system/mysql.service.d/slice.conf
# vim: ft=systemd
[Service]
Slice=contestant.slice
