# packer/files/itamae/roles/ci/files/etc/systemd/system/benchmarker.slice
# vim: ft=systemd
[Unit]
Description=benchmarker slice
Before=slices.target

[Slice]
Slice=isucon.slice
MemoryAccounting=yes
IOAccounting=yes
CPUAccounting=yes

MemoryHigh=2048M

CPUWeight=50
