# packer/files/itamae/roles/ci/files/etc/systemd/system/isucon.slice
# vim: ft=systemd
[Unit]
Description=isucon slice
Before=slices.target

[Slice]
MemoryAccounting=yes
IOAccounting=yes
CPUAccounting=yes
